# deploy/packtpub.timer
[Unit]
Description=Run packtpub

[Timer]
OnCalendar=07:50

[Install]
WantedBy=default.target
